# HelloSwiftWebWorld.conf
description "Hello Swift Web World"

start on startup

exec /home/ross/HelloSwiftWebWorld/.build/release/HelloSwiftWebWorld port=80 --workDir=/home/ross/HelloSwiftWebWorld
